# deploy/hubmq.service
[Unit]
Description=HubMQ — Unified communication hub
Documentation=http://localhost:3000/motreffs/hubmq
After=network-online.target nats.service
Requires=nats.service
Wants=network-online.target
# B1 — fallback on repeated failure (OnFailure appartient à [Unit], pas [Service])
OnFailure=hubmq-fallback-p0.service

[Service]
Type=simple
User=hubmq
Group=hubmq
Environment=RUST_LOG=info,hubmq_core=debug
Environment=HUBMQ_CONFIG=/etc/hubmq/config.toml
Environment=HUBMQ_CONF_D=/etc/hubmq/conf.d

# D3 — systemd LoadCredential (files remain root:root 600, service sees them via $CREDENTIALS_DIRECTORY)
LoadCredential=gmail-app-password:/etc/hubmq/credentials/gmail-app-password
LoadCredential=telegram-bot-token:/etc/hubmq/credentials/telegram-bot-token
LoadCredential=msg-relay-token:/etc/hubmq/credentials/msg-relay-token

ExecStart=/usr/local/bin/hubmq
# SIGHUP déclenche un reload du registry conf.d (hot-reload bots/listeners sans restart)
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Hardening
ProtectSystem=strict
# conf.d et credentials doivent être accessibles en écriture par le daemon
# pour que l'admin consumer puisse créer les tokens et modifier bots.toml.
ReadWritePaths=/var/lib/hubmq /var/log/hubmq /etc/hubmq/conf.d /etc/hubmq/credentials
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
