# deploy/agent/hubmq-agent-listener.service
[Unit]
Description=HubMQ Agent Listener — spawn claude-hubmq on Telegram messages
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=motreffs
Group=motreffs
WorkingDirectory=/home/motreffs
Environment=HOME=/home/motreffs
Environment=PATH=/home/motreffs/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/motreffs/.hubmq-agent/wrapper/hubmq-agent-listener.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
